# daemon/service/padkitd.service
# PadKit companion daemon — systemd USER unit
# SPDX-License-Identifier: MIT
#
# Install (per-user, no root needed for the service itself):
#   mkdir -p ~/.config/systemd/user
#   cp padkitd.service ~/.config/systemd/user/
#   # edit ExecStart below to the absolute path of your padkitd binary
#   systemctl --user daemon-reload
#   systemctl --user enable --now padkitd.service
#   systemctl --user status padkitd.service
#   journalctl --user -u padkitd -f
#
# HID access on Linux needs the udev rule (service/99-padkit.rules) installed
# once as root so your user can open the vendor interface without sudo.
#
# To have the user service start at boot (not just at login):
#   sudo loginctl enable-linger "$USER"

[Unit]
Description=PadKit companion daemon
Documentation=https://github.com/padkit/padkit
After=graphical-session.target

[Service]
Type=simple
# Absolute path required. Adjust to where you installed the binary:
ExecStart=%h/.local/bin/padkitd serve
Restart=on-failure
RestartSec=2
# The daemon writes its config here on first run if missing:
Environment=PADKIT_CONFIG=%h/.config/padkit/config.yaml

[Install]
WantedBy=default.target
